# tools_docker/Debian_Uwsgi/new_service.service.sample
[Unit]
Description=
After=network.target

[Service]
WorkingDirectory=/tracim/
Type=simple
RemainAfterExit=no
ExecStart=

[Install]
WantedBy=multi-user.target